# FastAPIApp.service
[Unit]
Description=Generic FastAPI Application Service
After=network.target

[Service]
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/venv/bin/uvicorn __APP_MODULE__ --host 0.0.0.0 --port __APP_PORT__
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=10
SyslogIdentifier=__PROJECT_NAME__

[Install]
WantedBy=multi-user.target


